# egnyte-sync.service
## Example systemd unit for Egnyte sync service
[Unit]
Description=Egnyte Sync Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/secho/Develop/Egnyte
ExecStart=/home/secho/Develop/Egnyte/.venv/bin/egnyte-cli service run --remote-interval 15
User=secho
Environment=HOME=/home/secho
Environment=XDG_CONFIG_HOME=/home/secho/.config
Environment=PYTHONUNBUFFERED=1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
